# deploy/systemd/chatgpt-connector-rollback.service
[Unit]
Description=Process a queued ChatGPT Connector rollback
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
ExecStart=/app/chatgpt_connector/deploy/process-rollback.sh
