# deploy/systemd/chatgpt-connector-health.timer
[Unit]
Description=Periodic ChatGPT Connector health check

[Timer]
OnBootSec=2m
OnUnitActiveSec=5m
AccuracySec=30s
Unit=chatgpt-connector-health.service

[Install]
WantedBy=timers.target
